# bbb.conf
[program:bbb]
command=python3 /home/curator2/bbb/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/bbb/test.err.log
stdout_logfile=/home/curator2/bbb/test.out.log